# infra/systemd/chickens-watcher.path
[Path]
PathModified=/srv/chickens

[Install]
WantedBy=multi-user.target