# infra/systemd/chickens-watcher.path
[Path]
PathModified=/srv/chickens/chicken-master.jar

[Install]
WantedBy=multi-user.target